# cloudflared.service
[Unit]
Description=Cloudflare Tunnel for RepublicAI Compute
After=network-online.target republic-http.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --config /etc/cloudflared/config.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
